# services/nginx/run
#!/bin/sh
cd ~/code/network-pinger/
exec nginx -c nginx/nginx.conf
