# services/orbited/run
#!/bin/sh
cd ~/code/network-pinger/
exec orbited --config nginx/orbited.cfg
