# services/paster/run
#!/bin/sh
cd ~/code/network-pinger/
exec paster serve production.ini --reload
