# services/paster/run
#!/bin/sh
exec 2>&1
cd ~/code/network-pinger/
exec paster serve production.ini --reload
